# stack-doctor.service.example
# Host-mode: run stack-doctor as a systemd service ON the host next to decypharr,
# so it has native power (restart decypharr locally, read its journal, touch the library)
# with no container-to-host bridge / SSH key needed.
#
#   cp doctor.py /opt/stack-doctor/doctor.py
#   cp stack-doctor.service.example /etc/systemd/system/stack-doctor.service   # edit values
#   systemctl daemon-reload && systemctl enable --now stack-doctor.service
#   journalctl -u stack-doctor -f

[Unit]
Description=stack-doctor - media stack health daemon (host mode)
After=network-online.target decypharr.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/stack-doctor/doctor.py
Restart=always
RestartSec=30
# --- mode ---
Environment=DOCTOR_MODE=cron
Environment=DOCTOR_INTERVAL=600
Environment=DOCTOR_LOG_FILE=/var/log/stack-doctor.log
# --- checks (host can run them all) ---
Environment=ENABLE_QUEUE=true
Environment=ENABLE_PROVIDERS=true
Environment=ENABLE_DECYPHARR=true
Environment=ENABLE_PLEX=true
Environment=ENABLE_RESOURCES=true
Environment=ENABLE_JANITOR=true
Environment=ENABLE_WARMER=true
Environment=ENABLE_UI=true
Environment=DOCTOR_UI_PORT=12345
Environment=DOCTOR_BLOCKLIST=true
# churn brake: after 3 dead grabs of the same episode/movie, un-monitor it, then retry in 24h
Environment=DOCTOR_CHURN_LIMIT=3
Environment=DOCTOR_CHURN_ACTION=backoff
Environment=DOCTOR_CHURN_COOLDOWN=86400
# --- *arr instances (reach over the LAN); add prowlarr for the providers check ---
Environment=INSTANCE_1_NAME=sonarr   INSTANCE_1_TYPE=sonarr    INSTANCE_1_URL=http://ARR_HOST:8989  INSTANCE_1_APIKEY=xxxx
Environment=INSTANCE_2_NAME=radarr   INSTANCE_2_TYPE=radarr    INSTANCE_2_URL=http://ARR_HOST:7878  INSTANCE_2_APIKEY=xxxx
Environment=INSTANCE_3_NAME=prowlarr INSTANCE_3_TYPE=prowlarr  INSTANCE_3_URL=http://ARR_HOST:9696  INSTANCE_3_APIKEY=xxxx
# --- decypharr (LOCAL: native restart + journal, no ssh) ---
Environment=DECYPHARR_MOUNT_TEST=/mnt/library/movies
Environment=DECYPHARR_RESTART_CMD=systemctl restart decypharr
Environment=DECYPHARR_URL=http://127.0.0.1:8282
# --- plex (also used by the warmer) ---
Environment=PLEX_URL=http://PLEX_HOST:32400
Environment=PLEX_TOKEN=xxxx
# --- warmer (precache likely-next media so playback starts instantly) ---
Environment=WARMER_PRECACHE_MB=64
Environment=WARMER_TAIL_MB=8
Environment=WARMER_SOURCES=ondeck,next
Environment=WARMER_MAX_PER_CYCLE=40
Environment=WARMER_LOAD_MAX=30
# detail-page warming: warm the exact title a viewer opens (tail Plex's server log).
# host can reach Plex-in-a-Proxmox-container with: pct exec <ctid> -- tail -n0 -F '<logpath>'
Environment=WARMER_PLEXLOG_CMD=pct exec 100 -- tail -n0 -F '/mnt/appdata/plex/Library/Application Support/Plex Media Server/Logs/Plex Media Server.log'
# --- janitor reads decypharr's journal directly ---
Environment=JANITOR_LIBRARY_PATHS=/mnt/library/tv,/mnt/library/movies
Environment=JANITOR_LOG_CMD=journalctl -u decypharr -n 10000 --no-hostname
Environment=JANITOR_QUARANTINE_DIR=/mnt/library/.quarantine

[Install]
WantedBy=multi-user.target
